# deployment/systemd/bitaxe-autopilot.service
[Unit]
Description=BitAxe Autopilot Optimization Service V2.0.0
After=network.target bitaxe-web.service
Wants=network.target

[Service]
Type=simple
User=bitaxe
Group=bitaxe
WorkingDirectory=/opt/bitaxe-web
Environment=PATH=/opt/bitaxe-web/venv/bin
Environment=DATABASE_PATH=/opt/bitaxe-web/data/bitaxe_data.db
ExecStart=/opt/bitaxe-web/venv/bin/python scripts/autopilot_service_runner.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target